# laptop/.config/upstart/desktopClose.conf
description "Desktop Close Task"
start on session-end
task
script
mpc -p 6601 pause >> /dev/null
# Commands to unmount sshfs go here as shell-script.
end script
