# laptop/.config/upstart/desktopOpen.conf
description "Desktop Open Task"
start on desktop-start
task
script
# Commands to mount sshfs go here as shell-script.
mpc -p 6601 play >> /dev/null
setxkbmap -option caps:super

end script
